# debian/lxcfs.upstart
description "FUSE filesystem for LXC"
author "Stéphane Graber <stgraber@ubuntu.com>"

respawn

start on starting lxc or cgmanager-ready
stop on runlevel [06]

exec /usr/bin/lxcfs -s -f -o allow_other /var/lib/lxcfs
